# services/setup-traefik.service
[Service]
Type=simple
ExecStart=/bin/sh /etc/setup-traefik.sh

[Install]
WantedBy=multi-user.target